# supervisor.conf
; ==============================================================================
; Supervisor Configuration for Laravel Inventory Management System
; Server Path: /home/u776021627/domains/inventory.satmec.com/public_html
; ==============================================================================
;
; INSTALLATION STEPS:
;
; 1. SSH into your server
;
; 2. Check if supervisor is installed:
;      which supervisord
;    If not installed:
;      sudo apt-get update && sudo apt-get install supervisor
;
; 3. Copy this config to supervisor:
;      sudo cp /home/u776021627/domains/inventory.satmec.com/public_html/supervisor.conf /etc/supervisor/conf.d/inventory-system.conf
;
; 4. Reload supervisor:
;      sudo supervisorctl reread
;      sudo supervisorctl update
;      sudo supervisorctl start inventory-system:*
;
; 5. Check status:
;      sudo supervisorctl status
;
; USEFUL COMMANDS:
;   sudo supervisorctl status                      - Check status of all processes
;   sudo supervisorctl restart inventory-system:*  - Restart all processes
;   sudo supervisorctl stop inventory-system:*     - Stop all processes
;   sudo supervisorctl start inventory-system:*    - Start all processes
;   sudo supervisorctl tail -f inventory-worker-00 - Watch worker logs
;   sudo supervisorctl tail -f inventory-scheduler - Watch scheduler logs
;
; ==============================================================================

; ------------------------------------------------------------------------------
; Laravel Queue Worker
; Processes background jobs from the queue
; ------------------------------------------------------------------------------
[program:inventory-worker]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /home/u776021627/domains/inventory.satmec.com/public_html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --memory=512
directory=/home/u776021627/domains/inventory.satmec.com/public_html
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=u776021627
numprocs=2
redirect_stderr=true
stdout_logfile=/home/u776021627/domains/inventory.satmec.com/public_html/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=3600

; ------------------------------------------------------------------------------
; Laravel Scheduler
; Runs scheduled tasks every minute (token refresh, delivery checks, etc.)
; ------------------------------------------------------------------------------
[program:inventory-scheduler]
process_name=%(program_name)s
command=/bin/bash -c "while true; do /usr/bin/php /home/u776021627/domains/inventory.satmec.com/public_html/artisan schedule:run --verbose --no-interaction; sleep 60; done"
directory=/home/u776021627/domains/inventory.satmec.com/public_html
autostart=true
autorestart=true
user=u776021627
numprocs=1
redirect_stderr=true
stdout_logfile=/home/u776021627/domains/inventory.satmec.com/public_html/storage/logs/scheduler.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

; ------------------------------------------------------------------------------
; Group Configuration
; Allows controlling all processes together
; ------------------------------------------------------------------------------
[group:inventory-system]
programs=inventory-worker,inventory-scheduler
priority=999
